# configs/etc/cron.d/server.info
@reboot dmidecode | txt2html > /var/www/html/index.html



